# packaging/linux/expense-tracker-keycloak.service
[Unit]
Description=Expense Tracker Keycloak Identity Server
Documentation=https://www.keycloak.org/documentation
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=simple
User=expense-tracker
Group=expense-tracker
WorkingDirectory=/opt/expense-tracker/keycloak-server
Environment=KC_HTTP_PORT=8080
Environment=KC_DB=postgres
Environment=KC_HEALTH_ENABLED=true
EnvironmentFile=-/etc/expense-tracker/keycloak.conf
ExecStart=/opt/expense-tracker/keycloak-server/bin/kc.sh start --optimized
Restart=on-failure
RestartSec=15
StandardOutput=append:/var/log/expense-tracker/keycloak.log
StandardError=append:/var/log/expense-tracker/keycloak-error.log

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ReadWritePaths=/var/log/expense-tracker /opt/expense-tracker/keycloak-server/data

[Install]
WantedBy=multi-user.target
